# systemd/jira-robot-sync.timer
[Unit]
Description=Run the Jira robot sync nightly

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
