# launchd/com.decalscanner.serve-robots.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.decalscanner.serve-robots</string>

  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/decal-scanner/scripts/venv/bin/python</string>
    <string>__HOME__/decal-scanner/scripts/serve_robots.py</string>
  </array>

  <!-- Runs continuously: started at load and restarted if it ever exits. -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/decal-scanner-serve.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/decal-scanner-serve.log</string>
</dict>
</plist>
